# .gc/scripts/systemd/gc-dolt-maintenance.service
[Unit]
Description=Gas City Dolt Maintenance (gc, cleanup)

[Service]
Type=oneshot
WorkingDirectory=/mnt/ext-fast/gc
ExecStart=/mnt/ext-fast/gc/scripts/dolt-maintenance.sh
Environment=HOME=/home/asiri
Environment=PATH=/home/asiri/go/bin:/home/asiri/.local/bin:/usr/local/bin:/usr/bin:/bin
