# .gc/scripts/systemd/gc-health.service
[Unit]
Description=Gas City Health Check

[Service]
Type=oneshot
WorkingDirectory=/mnt/ext-fast/gc
ExecStart=/usr/bin/python3 /mnt/ext-fast/gc/scripts/health.py --fix
Environment=HOME=/home/asiri
Environment=PATH=/home/asiri/go/bin:/home/asiri/.local/bin:/usr/local/bin:/usr/bin:/bin
